# systemd/sshd-report@.timer
[Unit]
Description=Send %i SSHD email reports

[Timer]
OnCalendar=%i
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
